# build/supervisord.conf
[supervisord]
loglevel=info
nodaemon=true
user=root
pidfile=/tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=run-nginx.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:bob]
command=run-bob.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
